# services/celery.service
[Unit]
Description=Celery Service
After=network.target redis-server.service
Requires=redis-server.service

[Service]
EnvironmentFile=/home/ubuntu/aigle-api/.env
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aigle-api
ExecStart=/home/ubuntu/aigle-api/venv/bin/celery -A aigle worker --loglevel=info --concurrency=1 --prefetch-multiplier=1 -Q celery,sequential_commands
Restart=always
RestartSec=10

# Warm shutdown: finish current task, refuse new ones. KillMode=mixed sends SIGTERM
# to the master only — the control-group default would also hit the prefork child
# and kill the task. Wait up to 12h for the task before SIGKILL.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=12h

[Install]
WantedBy=multi-user.target
